# status/crontab
*/5 * * * * python3 poll.py /var/www/html/js/stats.js